# public/assets/client_installer/payload/Library/LaunchDaemons/com.github.munkireport.runner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>KeepAlive</key>
	<dict>
		<key>PathState</key>
		<dict>
			<key>/Users/Shared/.com.github.munkireport.run</key>
			<true/>
		</dict>
	</dict>
	<key>Label</key>
	<string>com.github.munkireport.runner</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/munkireport/munkireport-runner</string>
		<string>--launchd</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StartInterval</key>
	<integer>3600</integer>
</dict>
</plist>
